# mkosi.skeleton/etc/systemd/system/nvidia-fabricmanager.service.d/override.conf
[Service]
# Create /run/nvidia-fabricmanager directory for PID file
RuntimeDirectory=nvidia-fabricmanager
